# crates/l2/tee/mkosi.extra/etc/systemd/system/payload.service
[Unit]
After=network-online.target

[Service]
ExecStart=/payload
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
